# devstack/service/openstack-plugin-flask.service
[Unit]
Description=Openstack Plugin Flask
After=network.target

[Service]
User=stack
Type=simple
WorkingDirectory=/opt/stack/openstack-plugin-flask
ExecStart=/opt/stack/openstack-plugin-flask/app/venv/bin/python3 /opt/stack/openstack-plugin-flask/app/app.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
